# systemd/clawos-watchdog.service
[Unit]
Description=ClawOS gateway watchdog
After=openclaw-gateway.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/openclaw-healthcheck.sh
