# image/overlays/etc/systemd/system/ssh.service.d/runtime-directory.conf
[Service]
RuntimeDirectory=sshd
RuntimeDirectoryMode=0755
